# rc.local.example
webhook -hooks /home/Giru/hooks.json &
python3 -m giru &

exit 0